# utils/autostart/squareline.desktop
[Desktop Entry]
Type=Application
Name=SquareLine Project
Exec=/home/pi/Desktop/x64-Debug/SquareLine_Project
